# templates/default/sidekiq.service.erb
[Unit]
Description=systemd script for running sidekiq for baritoMarket

[Service]
Restart=on-failure
RestartSec=2
StartLimitInterval=50
StartLimitBurst=10
ExecStart=/bin/bash -lc 'bundle exec sidekiq'
User=baritoMarket
WorkingDirectory=<%= @app_directory %>

[Install]
WantedBy=multi-user.target
